# src/release.timer
[Unit]
Description=Release Automat every week

[Timer]
OnCalendar=weekly
Persistent=true
Unit=release.service

[Install]
WantedBy=timers.target
